# amk-heritage-backend.service
[Unit]
Description=AMK Heritage Trail - FastAPI Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/amk_heritage/backend
ExecStart=/opt/amk_heritage/backend/venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
